# documents.service
[Unit]
Description=Document & Image Service (FastAPI + Gunicorn + Uvicorn workers)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data

# Must match the git checkout on the server
WorkingDirectory=/var/www/html/documents

# Leading "-" makes a missing .env non-fatal (API key is hardcoded in the app)
EnvironmentFile=-/var/www/html/documents/.env

Environment=PATH=/var/www/html/documents/venv/bin:/usr/local/bin:/usr/bin
Environment=PYTHONUNBUFFERED=1

ExecStart=/var/www/html/documents/venv/bin/gunicorn \
    -k uvicorn.workers.UvicornWorker \
    -w 4 \
    -b 0.0.0.0:8000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --access-logfile - \
    --error-logfile - \
    --capture-output \
    app.main:app

Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/html/documents/app/storage /var/www/html/documents/logs

LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
